# scripts/launchd/ai.domovina.apple-secret-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: mjesečno provjeri Apple "Sign in with Apple" client secret i
  rotiraj ga ako mu istječe za < THRESHOLD_DAYS (default 30). Skripta sama
  odlučuje hoće li ista rotirati — ako je iznad thresholda, izađe bez promjena.

  Instalacija + objašnjenje: docs/apple-secret-rotation.md
  Ovo je TEMPLATE u repou; instalirana kopija živi u
  ~/Library/LaunchAgents/ai.domovina.apple-secret-rotate.plist

  NAPOMENA: NODE_BIN je pinned na trenutni nvm node. Ako nadogradiš node
  (nvm install/use), ažuriraj NODE_BIN ovdje i u instaliranoj kopiji.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.domovina.apple-secret-rotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ms/git/domovinatv/domovina-api/scripts/rotate-apple-secret.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ms/git/domovinatv/domovina-api</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_BIN</key>
        <string>/Users/ms/.nvm/versions/node/v24.16.0/bin/node</string>
        <key>THRESHOLD_DAYS</key>
        <string>30</string>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin</string>
    </dict>

    <!-- 1. u mjesecu u 04:00 lokalno (Europe/Zagreb). Ako je Mac uspavan,
         launchd pokrene job pri sljedećem buđenju. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ms/git/domovinatv/domovina-api/logs/apple-secret-rotate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ms/git/domovinatv/domovina-api/logs/apple-secret-rotate.log</string>
</dict>
</plist>
